# net-wireless/rtl8192fu-driver/files/99-rtl8192fu-dongle.rules
ACTION=="add", DRIVERS=="sr", ATTRS{idVendor}=="0bda", RUN+="/usr/bin/eject /dev/$name & echo $name > /tmp/rtl.info"
