# net-wireless/rtl8192fu-driver/files/99-rtl8192fu-dongle.rules
ACTION=="add", DRIVER=="usb", ATTR{idVendor}=="0bda", ATTR{idProduct}=="a192", RUN+="/lib/udev/eject.sh"
