# net-misc/quectel-cm/files/init-quectel-service.conf
# Copyright 2022 The FydeOS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "FydeOS device quectel service."
author          "yang@fydeos.io"

start on started boot-services
stop on stopping boot-services

respawn
oom score 0

exec /usr/bin/quectel-CM
